# config/supervisor/drt.ini
[program:drt]
command=/home/your-username/.local/bin/gunicorn es.wsgi -b 0.0.0.0:8000
directory=/home/your-username/app/
user=your-username
autorestart=true
stdout_logfile = /home/your-username/gunicorn_supervisor.log
redirect_stderr=true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8